# auto.cron

*/6 */6 */1 * * python /Users/navjotbabrah/Documents/boost/lumi.py
#≈*/1 * * * * python /Users/navjotbabrah/Documents/boost/autoround.py
